# systemd/system/preston-brit-bisque.timer
[Unit]
Description=Trigger BRIT Bisque tracker

[Timer]
OnCalendar=monthly
RandomizedDelaySec=86400

[Install] 
WantedBy= timers.target
