# systemd/system/preston-checklistbank.timer
[Unit]
Description=Run Preston DataONE

[Timer]
OnCalendar=monthly
RandomizedDelaySec=86400

[Install] 
WantedBy= timers.target
