# systemd/system/preston-brit.timer
[Unit]
Description=Trigger BRIT tracker

[Timer]
OnCalendar=monthly
RandomizedDelaySec=86400

[Install] 
WantedBy= timers.target
