# systemd/system/preston-wikidata.timer
[Unit]
Description=Trigger Wikidata tracker

[Timer]
OnCalendar=monthly
RandomizedDelaySec=86400

[Install] 
WantedBy= timers.target
