# systemd/system/preston-bionomia.timer
[Unit]
Description=Trigger Bionomia tracker

[Timer]
OnCalendar=monthly
RandomizedDelaySec=86400

[Install] 
WantedBy=timers.target
