# systemd/system/preston.timer
[Unit]
Description=Run Preston

[Timer]
OnCalendar=monthly
RandomizedDelaySec=86400

[Install] 
WantedBy= timers.target
