# systemd/system/preston-obis.timer
[Unit]
Description=Run Preston OBIS

[Timer]
OnCalendar=monthly
RandomizedDelaySec=86400

[Install]
WantedBy= multi-user.target
