# systemd/system/preston-catalogueoflife.timer
[Unit]
Description=Run Preston Catalogue of Life tracker

[Timer]
OnCalendar=monthly
RandomizedDelaySec=86400

[Install] 
WantedBy= timers.target
